# darwin/Library/LaunchAgents/local.startup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.startup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>/Users/iam/.local/share/ilaumjd/startup.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/iam/Library/Logs/local.startup.out</string>

    <key>StandardErrorPath</key>
    <string>/Users/iam/Library/Logs/local.startup.err</string>
</dict>
</plist>
